# app/gentoo/local.d/save_screen_brightness.stop
#!/bin/sh

cat /sys/class/backlight/intel_backlight/brightness > /etc/local.d/prev_brightness
