# archive/gentoo/local.d/recover_screen_brightness.start
#!/bin/sh

cat /etc/local.d/prev_brightness > /sys/class/backlight/intel_backlight/brightness
